# install/systemd/chocopi.service
[Unit]
Description=Choco Voice Assistant
After=network.target bluetooth.service
After=pipewire.service pipewire-pulse.service wireplumber.service
Wants=network.target

[Service]
Type=simple
User=chocopi
Group=chocopi
WorkingDirectory=/opt/chocopi
Environment=PATH=/opt/chocopi/.venv/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/chocopi/chocopi
Restart=always
RestartSec=10
TimeoutStopSec=5
KillMode=mixed
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
